# root/base/projectroot/usr/lib/systemd/system/app.service
[Unit]
Wants=systemd-networkd-wait-online@sys.service
After=systemd-networkd-wait-online@sys.service

[Service]
SyslogIdentifier=app
ExecStart=/usr/lib/app/app
