# root/projectroot/lib/systemd/system/usbipd.service
[Service]
ExecStartPre=/sbin/modprobe usbip-host
ExecStart=/usr/sbin/usbipd
