# root/projectroot/usr/lib/systemd/system/gui.service
[Service]
WorkingDirectory=/usr/lib/gui
ExecStart=/usr/bin/python -B -u /usr/lib/gui/gui
