# root/projectroot/usr/lib/systemd/system/gui.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/bin/python -B -u /usr/lib/gui/server

Restart=always
RestartSec=1
StartLimitInterval=0
