# root/projectroot/usr/lib/systemd/system/studio.socket
[Socket]
ListenStream=8000
BindToDevice=sys

[Install]
WantedBy=sockets.target
