# root/projectroot/etc/systemd/system/debug.target
[Unit]
StopWhenUnneeded=true

Wants=mount-root-rw.service
