# root/projectroot/etc/systemd/system/mount-root-rw.service
[Unit]
StopWhenUnneeded=true

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/mount -o remount,rw /
ExecStop=/bin/mount -o remount,ro /
