# root/projectroot/etc/systemd/system/update-apply.service
[Unit]
ConditionPathExists=/mnt/init/update

[Service]
Type=oneshot
ExecStartPre=-/bin/sleep 3
ExecStart=/bin/systemctl --no-block reboot
